# config/rootfs/overlay/etc/systemd/system/thermal-mgmt.service
[Unit]
Description=EdgeNOS Thermal Management
After=platform-init.service
Requires=platform-init.service

[Service]
Type=simple
ExecStart=/usr/sbin/thermal-mgmt.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
